# scripts/za.co.landcruisersa.ratings.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>za.co.landcruisersa.ratings</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>/Users/wesleyroos/Developer/LandCruiserSA/scripts/refresh-ratings-cron.sh</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key><integer>1</integer>
    <key>Hour</key><integer>7</integer>
    <key>Minute</key><integer>0</integer>
  </dict>
  <key>RunAtLoad</key>
  <false/>
  <key>StandardOutPath</key>
  <string>/Users/wesleyroos/Library/Logs/lcsa-ratings.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/wesleyroos/Library/Logs/lcsa-ratings.log</string>
</dict>
</plist>
